# zl/uwsgi.conf
# uwsgi - uWSGI server

description "uWSGI server"

start on socket PROTO=inet PORT=8001
stop on runlevel [!2345]

exec uwsgi --ini /var/spool/www/zl/uwsgi_app.ini
